# files/pkg-plist
/usr/local/bin/helloworld
/usr/local/etc/rc.d/helloworld
